# systemd/.config/systemd/user/gphotos-sync.timer
[Unit]
Description=Run gphotos-sync hourly and on boot

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min 

[Install]
WantedBy=timers.target
